# deploy/supervisor-worker.conf
; BlastIQ queue worker — Supervisor program definition.
;
; This is the production-recommended way to run the queue worker:
; supervisor restarts it on crash, on reboot, and after every deploy
; (via `supervisorctl restart blastiq-worker:*` in deploy.sh).
;
; Why queue:work instead of Horizon:
;   Horizon depends on the predis/predis composer package OR phpredis
;   PHP extension. If your config/queue.php uses redis as the connection
;   AND neither is installed, `php artisan horizon` errors with
;   "Class Predis\Client not found" before it ever starts. Plain
;   `queue:work` reads the same connection but does not require Predis
;   for the master process.
;
; To switch to Horizon later:
;   1. composer require predis/predis  (or install php-redis extension)
;   2. Edit this file: replace the command= line with
;        command=php __PROJECT_PATH__/artisan horizon
;   3. Set numprocs=1 (Horizon is a single master that spawns its own
;      child workers based on config/horizon.php).
;   4. supervisorctl reread && supervisorctl update
;
; Installation:
;   sudo bash deploy/install-supervisor.sh
;   The script substitutes __PROJECT_PATH__ and __RUN_AS_USER__ into
;   /etc/supervisord.d/blastiq-worker.ini, then reloads supervisor.

; Queue priority order (left = highest):
;   default  — system fallback, also for CampaignBatchDispatch fan-out
;   messages — per-message send (SendWhatsAppMessage), high volume during blasts
;   imports  — ProcessContactImport CSV/XLSX parse + insert; rare but blocks
;              the operator's success flash because the controller dispatches
;              and exits. If 'imports' isn't in this list, uploads "succeed"
;              and no contacts ever appear (the job sits in the jobs table).
[program:blastiq-worker]
process_name=%(program_name)s_%(process_num)02d
command=php __PROJECT_PATH__/artisan queue:work --queue=default,messages,imports --tries=3 --timeout=120 --sleep=3
directory=__PROJECT_PATH__
autostart=true
autorestart=true
user=__RUN_AS_USER__
numprocs=1
redirect_stderr=true
stdout_logfile=__PROJECT_PATH__/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
; Long stopwaitsecs so a worker mid-job (e.g., sending a campaign batch)
; finishes cleanly on `supervisorctl restart` instead of being SIGKILL'd.
stopwaitsecs=3600
