# deploy/supervisor-reverb.conf
; BlastIQ Reverb daemon — Supervisor program definition.
;
; Phase 17 dependency. Reverb is a long-lived WebSocket server bound to
; 127.0.0.1:8080. Public traffic enters via HTTPS at /app, nginx proxies
; the WebSocket upgrade to this daemon (see deploy/nginx.conf).
;
; Installation:
;   sudo bash deploy/install-reverb.sh
;
; The script substitutes __PROJECT_PATH__ and __RUN_AS_USER__ into
; /etc/supervisord.d/blastiq-reverb.ini, then reloads supervisor.

[program:blastiq-reverb]
process_name=%(program_name)s
command=php __PROJECT_PATH__/artisan reverb:start --host=127.0.0.1 --port=8080
directory=__PROJECT_PATH__
autostart=true
autorestart=true
user=__RUN_AS_USER__
redirect_stderr=true
stdout_logfile=__PROJECT_PATH__/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=10
